# quadlets/caddy.container
[Unit]
Description=caddy

[Service]
Restart=always
RuntimeDirectory=program-explorer
RuntimeDirectoryPreserve=yes
AllowedCPUs=0-3

[Container]
Image=localhost/pe-caddy:latest
Network=host
Volume=${RUNTIME_DIRECTORY}:${RUNTIME_DIRECTORY}:z
Volume=caddy-files.volume:/www
Volume=caddy-data.volume:/data
Environment=RUNTIME_DIRECTORY=${RUNTIME_DIRECTORY}

[Install]
WantedBy=default.target
